# whmanager/testdata/webhost/php-fpm/conf.d/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/run/supervisord.pid

[program:php-fpm9000]
user=root
command=php-fpm83 -y /etc/php83/pool.d/php-fpm9000.conf -F 
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=false
startretries=0

[program:php-fpm9001]
command=php-fpm83 -y /etc/php83/pool.d/php-fpm9001.conf -F 
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=false
startretries=0

[program:php-fpm9002]
command=php-fpm83 -y /etc/php83/pool.d/php-fpm9002.conf -F 
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=false
startretries=0
